# scripts/jenkins-custom_update.timer
echo "[Unit]
Description=Run script to update Jenkins version

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target"
